# projects/kubernetes-sigs/kind/images/node/files/etc/systemd/system/etcd.service.d/env.conf
[Service]
Environment="ETCD_UNSUPPORTED_ARCH=arm64"
